# ansible/playbooks/roles/backup/templates/backup.service.j2
[Unit]
Description=ASPHub 3-2-1 Backup Service
After=network.target docker.service plakar-agent.service
Wants=plakar-agent.service

[Service]
Type=oneshot
User={{ backup_user }}
Group={{ backup_user }}
WorkingDirectory=/opt/backup
ExecStart=/opt/backup/scripts/backup-orchestrator.sh
TimeoutStartSec=3600
EnvironmentFile=/opt/backup/.env

# Logging
StandardOutput=append:/var/log/backup/backup.log
StandardError=append:/var/log/backup/backup.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target 